# confile/etc/profile.d/pipenv.sh
#!/usr/bin/env bash

: <<!
Name: pipenv.sh
Author: YJ
Email: yj1516268@outlook.com
Created Time: 2020-12-09 15:08:21

Description: 配置pipenv变量
!

export PIPENV_PYPI_MIRROR="http://mirrors.aliyun.com/pypi/simple"
